# .bunnyshell/docker/php/files/crontab
* * * * * php /srv/sylius/bin/console sylius:remove-expired-carts
* * * * * php /srv/sylius/bin/console sylius:cancel-unpaid-orders
